# roles/os-hardening/templates/etc/profile.d/cis_controls.sh.j2
{{ ansible_managed | comment }}
readonly TMOUT={{ os_user_timeout }} ; export TMOUT
umask {{ os_env_umask }}

